# crontab
0 8,20 * * * /usr/local/bin/php /var/www/cron_src.php >> /var/www/cron.log
